# dottest.service
[Unit]
Description=MyApp Daemon Service
After=network.target

[Service]
ExecStart=/usr/bin/myapp daemon
Restart=on-failure
User=myappuser
Group=myappuser

[Install]
WantedBy=multi-user.target
